# recore.docker/profiles/generic/lib64/systemd/system/fleet.service.d/fleet-interface.conf
[Service]
Environment="FLEET_ETCD_SERVERS=http://172.17.42.1:4001"
