# recore.docker/profiles/bootserver-metal/lib64/systemd/system/btrfs.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
After=systemd-cryptsetup@enc.service
After=systemd-cryptsetup@enc2.service

Before=umount.target

[Mount]
What=/dev/mapper/enc
Where=/btrfs
Type=btrfs
Options=rw,device=/dev/mapper/enc,device=/dev/mapper/enc2