# recore.docker/profiles/bootserver/lib64/systemd/system/var.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
After=systemd-cryptsetup@enc.service
Before=umount.target

[Mount]
What=/dev/mapper/enc
Where=/var
Type=btrfs
Options=defaults,subvol=var