# recore.config/btrfs.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
After=systemd-cryptsetup@enc.service
Before=local-fs.target

[Mount]
What=/dev/mapper/enc
Where=/btrfs
Type=btrfs
Options=rw
